# Dockerfile.monitor.supervisord.conf
[program:gpufl-monitor]
command=/usr/local/bin/gpufl-monitor
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=GPUFL_MONITOR_APP="%(ENV_GPUFL_MONITOR_APP)s",GPUFL_MONITOR_LOG_DIR="%(ENV_GPUFL_MONITOR_LOG_DIR)s",GPUFL_MONITOR_INTERVAL_MS="%(ENV_GPUFL_MONITOR_INTERVAL_MS)s"

[program:gpufl-agent]
command=java -jar /opt/gpufl-agent/agent.jar
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; All GPUFL_* env vars are inherited from the container environment.
; Key variables (set via `docker run -e` or docker-compose environment:):
;   GPUFL_SOURCE_FOLDER  — must match the log dir used by gpufl-monitor
;   GPUFL_SOURCE_PREFIX  — must match GPUFL_MONITOR_LOG_DIR base name (default: session)
;   GPUFL_PUBLISHER_TYPE — http or kafka
;   GPUFL_HTTP_URL       — e.g. http://backend:8080/api/v1/events/
;   GPUFL_HTTP_TOKEN     — Bearer token
;   GPUFL_LOG_TYPES      — default: device,scope,system (override to restrict channels)
;   GPUFL_CURSOR_FILE    — default: ./cursor.json (override for persistence across restarts)
